# libvirt/assets/usr/lib/systemd/system/libvirtd.service.d/10-logs.conf
[Service]
#FailureAction=poweroff
StandardOutput=tty
